# src/etc/cron.d/entrypoint.sh
#!/bin/bash

/etc/init.d/nginx start

# Monitoring Load on the server
# while :
# do

# sleep 10
# command

# break
# done

